# ops/systemd/stream-platform-healthcheck.timer
[Unit]
Description=Run Stream Platform Health Summary

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
